# mentor-service/deploy/shellnotes-mentor.service
[Unit]
Description=Shellnotes Mentor BB
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/shellnotes-mentor
EnvironmentFile=/etc/shellnotes-mentor.env
Environment=HOME=/home/deploy
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/node /opt/shellnotes-mentor/server.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryMax=2G
TasksMax=64

[Install]
WantedBy=multi-user.target
